# contrib/com.jondanzig.susurrate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jondanzig.susurrate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/botclankman/.local/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/botclankman/code/susurrate</string>
        <string>susurrate</string>
        <string>serve</string>
        <string>--host</string>
        <string>127.0.0.1</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/botclankman/code/susurrate</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/botclankman/.local/bin</string>
        <key>SUSURRATE_MODEL</key>
        <string>/Users/botclankman/.local/share/susurrate/models/ggml-small.bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/botclankman/.local/share/susurrate/serve.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/botclankman/.local/share/susurrate/serve.log</string>
</dict>
</plist>
